# zapp/whisper.service
[Unit]
Description=Whisper HTTP Server (whisper.cpp)
After=network.target

[Service]
Type=simple
# User should be set to the actual username during deployment
# Example: User=<username> or User=your_username
User=<username>
ExecStart=${HOME}/bin/whisper-server-wrapper.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
